# debian/ksuperkey/etc/xdg/autostart/ksuperkey.desktop
[Desktop Entry]
Name=ksuperkey
Comment=Open the application launcher in KDE Plasma Desktop using the Super key.
Exec=ksuperkey
Type=Application